# build/ip_updater.service
[Unit]
Description=IP Updater Service
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/ip_updater -config=/etc/ip_updater/config.conf
Restart=always
RestartSec=10

# 超时设置
TimeoutStartSec=60
TimeoutStopSec=30
TimeoutSec=90

# 信号处理优化
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes

# 服务限制
LimitNOFILE=1024
MemoryMax=100M

[Install]
WantedBy=multi-user.target
